# get_wpa.service
[Unit]
Description=Auto-setup WPA

[Service]
Type=oneshot
ExecStart=/usr/bin/get_wpa

[Install]
WantedBy=multi-user.target
